# systemd/backup-db.service
[Unit]
Description=Film School Assistant - Database Backup
After=network.target

[Service]
Type=oneshot
User=oc_her
WorkingDirectory=/srv/openclaw-her/workspace/film-school-assistant
ExecStart=/bin/bash scripts/backup_db.sh
EnvironmentFile=/srv/openclaw-her/workspace/film-school-assistant/.env
NoNewPrivileges=true
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
