# systemd/cleanup-audio.service
[Unit]
Description=Film School Assistant - Audio Cleanup
After=network.target

[Service]
Type=oneshot
User=oc_her
WorkingDirectory=/srv/openclaw-her/workspace/film-school-assistant
ExecStart=/usr/bin/python3 scripts/cleanup_audio.py
EnvironmentFile=/srv/openclaw-her/workspace/film-school-assistant/.env
NoNewPrivileges=true
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
