# systemd/film-school-bot.service
[Unit]
Description=Film School Assistant Telegram Bot
After=network.target
StartLimitBurst=5

[Service]
Type=simple
User=youruser
WorkingDirectory=/path/to/repo
# Set ExecStart to the venv python3 inside your cloned repo, e.g.:
# ExecStart=/path/to/repo/.venv/bin/python3 src/bot.py
ExecStart=/path/to/repo/.venv/bin/python3 src/bot.py
Restart=on-failure
RestartSec=10
EnvironmentFile=/path/to/repo/.env
NoNewPrivileges=true
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
