# systemd/notify-failure@.service
[Unit]
Description=Notify Telegram on failure of %i
After=network.target

[Service]
Type=oneshot
User=oc_her
WorkingDirectory=/srv/openclaw-her/workspace/film-school-assistant
ExecStart=/usr/bin/python3 scripts/notify_failure.py
EnvironmentFile=/srv/openclaw-her/workspace/film-school-assistant/.env
Environment=FAILED_UNIT=%i
NoNewPrivileges=true
StandardOutput=journal
StandardError=journal
