# systemd/reminder.timer
[Unit]
Description=Film School Assistant — Daily Reminder Timer
Requires=reminder.service

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true

[Install]
WantedBy=timers.target
