# systemd/reminder.timer
[Unit]
Description=Film School Assistant — Daily Reminder Timer
Requires=reminder.service

[Timer]
OnCalendar=*-*-* *:00,15,30,45:00
Persistent=true

[Install]
WantedBy=timers.target
